# data/indicator-datetime.service.in
[Unit]
Description=Indicator Date & Time Backend
PartOf=graphical-session.target
PartOf=indicators.target

[Service]
ExecStart=@pkglibexecdir@/indicator-datetime-service
Restart=on-failure

[Install]
WantedBy=indicators.target
