# data/indicator-datetime.conf.in
description "Indicator Date & Time Backend"

start on indicators-loaded or indicator-services-start
stop on desktop-end or indicator-services-end

env G_MESSAGES_DEBUG=all
export G_MESSAGES_DEBUG

respawn

exec @libexecdir@/indicator-datetime-service
