# data/indicator-datetime.service.in
[Unit]
Description=Indicator Date & Time Backend
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=@pkglibexec@/indicator-datetime-service
Restart=on-failure
